# systemd/docstore-backup.service
# ═══════════════════════════════════════════════════════════════
# Agentischer Document Store — Backup-Timer
# Fuehrt alle 6 Stunden ein Backup aus
# ═══════════════════════════════════════════════════════════════
#
# Installation:
#   sudo cp systemd/docstore-backup.service /etc/systemd/system/
#   sudo cp systemd/docstore-backup.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now docstore-backup.timer
#
# Status:
#   sudo systemctl list-timers docstore-backup.timer
# ═══════════════════════════════════════════════════════════════

[Unit]
Description=Agentischer Document Store — Automatisches Backup
Requires=docstore.service
After=docstore.service

[Service]
Type=oneshot
WorkingDirectory=/opt/docstore-prod
EnvironmentFile=-/etc/default/docstore
ExecStart=/usr/bin/make backup
StandardOutput=journal
StandardError=journal
User=root

# Security
NoNewPrivileges=true
PrivateTmp=true
